# systemd/jarvis-voice-pipeline.service
[Unit]
Description=Jarvis Voice Pipeline (Telegram → faster-whisper → vault)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=lordmuffin
WorkingDirectory=/opt/jarvis-voice/src
EnvironmentFile=/home/lordmuffin/.jarvis.env
ExecStart=/home/lordmuffin/.whisper-venv/bin/python -m jarvis_voice.telegram_bot
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
