# systemd/mnt-notes.automount
[Unit]
Description=SSHFS vault automount
After=network-online.target

[Automount]
Where=/mnt/notes
TimeoutIdleSec=0

[Install]
WantedBy=multi-user.target
